# deploy/kindle-131-low-power/kindle-131-low-power-failsafe.conf
start on started lab126
task

exec /mnt/us/kindle-131-low-power-rollback-20260716-115926/failsafe-boot.sh --boot-check
